# david.conf
description "Starter script for helper bot David"
start on runlevel [2345]
stop on runlevel [!2345]

exec ~/Documents/personal/ai_and_data_science/trigger_word/alt.py